# mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=sh -c '--datadir='/var/lib/mysql' || kill 1'

[program:telegraf]
command=sh -c 'telegraf -config /telegraf/telegraf.toml || kill 1'