# influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=sh -c 'influxdb || kill 1'

[program:telegraf]
command=sh -c 'telegraf -config /telegraf/telegraf.toml || kill 1'